# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock
autostart=false
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac
autorestart=true
user=user
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pulseaudio]
command=pulseaudio --daemonize=no --disallow-exit --exit-idle-time=-1 --no-cpu-limit --log-target=stderr
autorestart=true
user=user
priority=12
environment=XDG_RUNTIME_DIR="/run/user/1001",HOME="/home/user"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xpra]
command=xpra start :99 --use-display=yes --bind-ws=0.0.0.0:10000 --no-daemon --no-notifications --start-new-commands=yes --html=on --pulseaudio=no --speaker=yes --microphone=no --input-method=none --webcam=no --printing=no
autorestart=true
user=user
priority=15
startsecs=3
environment=DISPLAY=":99",XDG_RUNTIME_DIR="/run/user/1001",HOME="/home/user"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server --bind 127.0.0.1 --port 6379 --save 60 1 --loglevel warning --dir /var/run/redis
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nextjs]
command=bun run start -p 3002
directory=/app
user=api
autorestart=true
priority=20
environment=NODE_ENV="production",IS_DOCKER="true",REDIS_HOST="127.0.0.1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ws-server]
command=node --import tsx src/server/ws-server.ts
directory=/app
user=api
autorestart=true
priority=20
environment=IS_DOCKER="true",REDIS_HOST="127.0.0.1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
